# backup/systemd/sysbackup_daily.timer
[Unit]
Description=Timer for daily sysbackup

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=multi-user.target

