# backup/systemd/sysbackup_monthly.timer
[Unit]
Description=Timer for monthly sysbackup

[Timer]
OnCalendar=*-*-1 04:30:00
Persistent=true

[Install]
WantedBy=multi-user.target

